# server/systemd/squidhallmu_restart.timer
[Timer]
OnCalendar=*-*-* 14:00:00 UTC

[Install]
WantedBy=timer.target

