# server/systemd/squidhallmu_restart.service
[Unit]
Description=SquidHallVR Colyseus server daily restart to clear the hall

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart squidhallmu
